# systemd/funding-alerts.service
[Unit]
Description=Funding Alerts Bot
After=network.target

[Service]
Type=simple
User=gpt
WorkingDirectory=/home/gpt/funding-alerts-bot
EnvironmentFile=/home/gpt/funding-alerts-bot/.env
ExecStart=/home/gpt/funding-alerts-bot/venv/bin/python main.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
